# trash/telegram_bot.service
[Unit]
Description=Football School Telegram Bot Manager (Multi-Tenancy)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/football_school
Environment="PATH=/opt/football_school/venv/bin"
Environment="APP_URL=https://d-promo.uz"
Environment="PYTHONUNBUFFERED=1"

# Запуск менеджера ботов (запускает отдельный бот для каждой школы)
ExecStart=/opt/football_school/venv/bin/python3 /opt/football_school/telegram_bot_manager.py

# Перезапуск при сбое
Restart=always
RestartSec=10

# Ограничения ресурсов
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target

